# nodepm/user/etc/autostart/nodepm-user.desktop
[Desktop Entry]
Type=Application
Categories=System
Name=NodePm User Program
Exec=bash -c "nodepm-user"
