# nodepm/0.5/main/etc/systemd/system/nodepm-init.service
[Unit]
Description=None-Desktop Power Manager (NODEPM) initialization service

[Service]
Type=simple
ExecStart=bash -ic "nodepm init"

[Install]
WantedBy=multi-user.target
