# nodepm-tray/pack/etc/xdg/autostart/nodepm-tray.desktop
[Desktop Entry]
Type=Application
Categories=GTK;System
Name=NodePm Tray
Exec=bash -c "nodepm-tray"
